# supervisor/consul-template.ini
[program:consultemplate]
command=/usr/local/bin/consul-template -consul=%(ENV_CONSUL_URL)s -template="/etc/varnish/consul.ctmpl:/etc/varnish/default.vcl:/usr/share/varnish/reload-vcl -q"
autorestart=true
